# systemd/.config/systemd/user/hypridle-watchdog.timer
[Unit]
Description=Check hypridle every 2 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=120

[Install]
WantedBy=timers.target
